# deployconfig/backupdb.crontab
0 */1 * * * python {remote_website_dir}/mysite/backupdb.py
